# kubeadm/config/rc.local
#!/usr/bin/env bash

/etc/sysctl.d
/etc/init.d/procps restart
